# Discord/deploy/vps/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel — Discord Unibox prod
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-networks/
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=notify
User=cloudflared
Group=cloudflared

# The config file is committed at deploy/cloudflare/tunnel-config.yaml.
# bootstrap.sh installs it to /etc/cloudflared/config.yml.
ExecStart=/usr/bin/cloudflared --no-autoupdate --config /etc/cloudflared/config.yml tunnel run
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s

# Sandboxing.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Logging.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
